# scripts/systemd/hippocampus-pull-latest.timer
[Unit]
Description=Hourly hippocampus weight refresh

[Timer]
OnCalendar=hourly
Persistent=true
AccuracySec=5m

[Install]
WantedBy=timers.target
