# scripts/systemd/thebrain-pi-daemon.service
[Unit]
Description=TheBrain Pi Daemon (replay writer)
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=image
WorkingDirectory=/home/image/thebrain
Environment=HIPPO_FRAME=/tmp/hippo_latest.jpg
Environment=HIPPO_SPEC=/tmp/hippo_spec.png
Environment=PI_REPLAY_DIR=/home/image/thebrain/replay
ExecStart=/home/image/thebrain/.venv/bin/python scripts/run_pi_daemon.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
